# Scripts/github-repos-scanner.service
[Unit]
Description=GitHub Repositories Scanner
After=network.target

[Service]
Type=oneshot
ExecStart=/home/duckyonquack999/.local/bin/github-repos-scanner.sh
User=duckyonquack999

[Install]
WantedBy=multi-user.target

